# templates/initscripts/firewall.init.ubuntu.v4.erb
#!/bin/sh
# 2014, j.grassler@syseleven.de

# Based on firewall.init.gentoo, adapted for Ubuntu's init peculiarities.

# Note: As of https://youtrack.syseleven.de/issue/pp-1222 both this file and
# firewall.init.gentoo need to be modified for changes to the vefirewall init
# script. 

### BEGIN INIT INFO
# Provides:          firewall
# Required-Start:    mountkernfs $local_fs
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# X-Start-Before:    $network
# X-Stop-After:      $network
# Short-Description: Set up iptables rules
### END INIT INFO

IPTABLES=iptables
no_rules=23
run_file=/run/firewall.running           # Records running state
<%= scope.function_template(["vefirewall/initscripts/firewall.init.ubuntu.erb"]) %>
